# oxd-server_1.0.3_all/DEBIAN/postinst
#!/bin/sh

set -e

update-rc.d oxd defaults
echo 'Starting oxd server...'
service oxd restart

exit 0
