# oxd-server_1.0.3_all/DEBIAN/prerm
#!/bin/sh

set -e
echo 'Stoping oxd server...'
service oxd stop
rm -f /opt/oxd-server/log/oxd*
rm -f /opt/oxd-server/pid
update-rc.d -f oxd remove

exit 0
